# data/eos-metrics-event-recorder.service.in
[Unit]
Description=EndlessOS Metrics Event Recorder Server
Requires=dbus.service
After=dbus.service

[Service]
User=metrics
Type=dbus
BusName=com.endlessm.Metrics
ExecStart=@libexecdir@/eos-metrics-event-recorder

[Install]
WantedBy=multi-user.target
