# data/eos-metrics-event-recorder.service.in
[Unit]
Description=EndlessOS Metrics Event Recorder Server
Requires=dbus.service
Wants=systemd-logind.service
After=dbus.service systemd-logind.service

[Service]
Environment=DCONF_PROFILE=/dev/null
User=metrics
Type=dbus
BusName=com.endlessm.Metrics
ExecStart=@libexecdir@/eos-metrics-event-recorder

[Install]
WantedBy=multi-user.target
